# pynq/spicer-daq.service
[Unit]
Description=Spicer DAQ Mode Manager (RUN or DEBUG mode)
After=network.target
# Start after all services are up
After=multi-user.target

[Service]
Type=simple
User=xilinx
Group=xilinx
WorkingDirectory=/var/www/html/pynq
ExecStartPre=/bin/sleep 30
ExecStart=/var/www/html/pynq/mode_manager.sh startup
Restart=always
RestartSec=10
Environment="PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"

[Install]
WantedBy=multi-user.target