# pynq/spicer-daq.service
[Unit]
Description=Spicer DAQ WebSocket Server
After=network.target

[Service]
Type=simple
User=xilinx
Group=xilinx
WorkingDirectory=/var/www/html/pynq
Environment="PATH=/usr/local/share/pynq-venv/bin:/usr/bin:/bin"
ExecStart=/usr/local/share/pynq-venv/bin/python /var/www/html/pynq/server.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target